# gateway/systemd/notify-gateway.service
[Unit]
Description=Notify Gateway — webhook dispatcher for Home Assistant notifications
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/hass-ai-gateway/gateway/services/notify_gateway
ExecStart=/root/hass-ai-gateway/gateway/services/notify_gateway/notify_gateway_env/bin/python /root/hass-ai-gateway/gateway/services/notify_gateway/main.py
Environment="PATH=/root/hass-ai-gateway/gateway/services/notify_gateway/notify_gateway_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
