# assistant/systemd/voice-gateway.service
[Unit]
Description=Voice Gateway for Home Assistant (RPi / ESP32 audio input)
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/smarthome/assistant/services/voice_gateway
ExecStart=/root/smarthome/assistant/services/voice_gateway/voice_gateway_env/bin/python /root/smarthome/assistant/services/voice_gateway/main.py
Environment="PATH=/root/smarthome/assistant/services/voice_gateway/voice_gateway_env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
