# addon/rootfs/etc/services.d/voice-gateway/run
#!/command/with-contenv bashio
set -e

if ! bashio::config.true 'enable_voice_gateway'; then
    bashio::log.info "[voice-gateway] disabled in config — sleeping"
    exec sleep infinity
fi

source /usr/lib/hass-ai-gateway/export-env.sh

bashio::log.info "[voice-gateway] starting on port ${GATEWAY_PORT}"
cd /opt/gateway/services/voice_gateway
exec python3 main.py
